# on_the_pi/crontab
# runs script once on reboot
@reboot /PATH/TO/SCRIPT/rtlC -d 24 &
